# snap/local/share/systemd/user/hermes-agent.service
[Unit]
Description=Hermes AI Agent gateway
After=network.target

[Service]
ExecStart=/snap/bin/hermes-agent.daemon
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
